# src/files/etc/supervisor/conf.d/default.conf
[program:php-fpm]
command=/usr/sbin/php-fpm%(ENV_PHP_VERSION)s --nodaemonize --fpm-config /etc/php/%(ENV_PHP_VERSION)s/fpm/php-fpm.conf
autostart=true
autorestart=true

[program:cron]
command=cron -f -L 15
autostart=true
autorestart=false
